# u80211d/files/etc/init.d/u80211d
#!/bin/sh /etc/rc.common
START=80

USE_PROCD=1
PROG=/usr/sbin/u80211d

start_service() {
	procd_open_instance
	procd_set_param command "$PROG"
	procd_close_instance
}
